# dni/go/scripts/service/dni.service.in
[Unit]
Description=DNI service

[Service]
Type=simple
ExecStart= /usr/local/dni/start_service.sh 

[Install]
WantedBy=multi-user.target